# projects/shared/lead-scraper/config/com.marceausolutions.sms-scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.marceausolutions.sms-scheduler</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>src.sms_scheduler</string>
        <string>run</string>
        <string>--window</string>
        <string>morning_peak</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/williammarceaujr./dev-sandbox/projects/lead-scraper</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <!-- Run at 9 AM EST on Tuesday through Friday -->
    <key>StartCalendarInterval</key>
    <array>
        <!-- Tuesday 9 AM -->
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- Wednesday 9 AM -->
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- Thursday 9 AM -->
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- Friday 9 AM -->
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/williammarceaujr./dev-sandbox/projects/lead-scraper/logs/sms_scheduler.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/williammarceaujr./dev-sandbox/projects/lead-scraper/logs/sms_scheduler_error.log</string>

    <!-- Don't run if already running -->
    <key>LegacyTimers</key>
    <true/>
</dict>
</plist>
