# scripts/launchd/com.marceau.bidirectional-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.marceau.bidirectional-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/williammarceaujr./dev-sandbox/scripts/bidirectional_sync.py</string>
        <string>--data-only</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/williammarceaujr./dev-sandbox</string>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>StandardOutPath</key>
    <string>/Users/williammarceaujr./dev-sandbox/projects/personal-assistant/logs/bidirectional-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/williammarceaujr./dev-sandbox/projects/personal-assistant/logs/bidirectional-sync-error.log</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
