# projects/shared-multi-tenant/lead-scraper/launchd/com.marceausolutions.campaign-launcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.marceausolutions.campaign-launcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>src.campaign_auto_launcher</string>
        <string>--for-real</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/williammarceaujr./dev-sandbox/projects/shared-multi-tenant/lead-scraper</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- SCHEDULES: Monday, Wednesday, Friday (HVAC) + Tuesday, Thursday, Saturday (Marceau) -->
    <key>StartCalendarInterval</key>
    <array>
        <!-- Monday 8 AM - HVAC -->
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>

        <!-- Tuesday 8 AM - Marceau Solutions -->
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>

        <!-- Wednesday 8 AM - HVAC -->
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>

        <!-- Thursday 8 AM - Marceau Solutions -->
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>

        <!-- Friday 8 AM - HVAC -->
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>

        <!-- Saturday 10 AM - Marceau Solutions (lighter batch) -->
        <dict>
            <key>Weekday</key>
            <integer>6</integer>
            <key>Hour</key>
            <integer>10</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/williammarceaujr./dev-sandbox/projects/shared-multi-tenant/lead-scraper/output/logs/campaign-launcher.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/williammarceaujr./dev-sandbox/projects/shared-multi-tenant/lead-scraper/output/logs/campaign-launcher-error.log</string>

    <!-- Don't run on load (only on schedule) -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
