# projects/shared/grok-agent/grok-agent.service
[Unit]
Description=Grok Agentic Assistant
After=network.target

[Service]
Type=simple
User=clawdbot
WorkingDirectory=/home/clawdbot/dev-sandbox/projects/shared/grok-agent
ExecStart=/usr/bin/python3.11 grok_agent_standalone.py
Restart=always
RestartSec=5
Environment=PATH=/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
